# setup/com.eightmen.ranking.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.eightmen.ranking</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/c5ysk6/Documents/eight-men-monthly-ranking/update.sh</string>
  </array>

  <!-- 毎日13時00分に実行 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>13</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/c5ysk6/Documents/eight-men-monthly-ranking/update.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/c5ysk6/Documents/eight-men-monthly-ranking/update.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
